# application/clamav/cron/s6/etc/services.d/clamd/finish
#!/usr/bin/execlineb -S0
backtick DATE { date -u --rfc-3339=ns }
importas date DATE

foreground { echo "[clamd] ${date} !! Clamd service exited for some reason. Stopping container." }
s6-svscanctl -t /var/run/s6/services  # exit all other services, kill the container
