# application/clamav/cron/s6/etc/services.d/clamd/run
#!/bin/execlineb -P
###!/usr/bin/with-contenv /bin/execlineb
#importas log_path LOG_PATH

s6-envdir -fn -- /var/run/s6/container_environment
backtick HOST_ID { hostname }
importas LOG_PATH LOG_PATH
importas HOST_ID HOST_ID

# Redirect stderr to stdout.
fdmove -c 2 1

# Notify service manager when service is up
s6-notifyoncheck -d -w 1000 -c "redirfd -w 1 /dev/null fdmove -c 2 1 stat /var/run/clamav/clamd.ctl"

redirfd -a 1 ${LOG_PATH}/logs-service/clamd_svc.${HOST_ID}.log
/usr/bin/nice -n 19 /usr/sbin/clamd --foreground --config-file=${LOG_PATH}/conf/clamd.conf
