# crond/s6/etc/services.d/cron/run
#!/bin/execlineb -P

# Redirect stderr to stdout.
fdmove -c 2 1

/usr/sbin/crond -f -c /etc/crontabs
